# hwoj.conf/hwoj-discordbot.conf
[program:hwoj-discordbot]
command=/home/oj/HWOJ/discordbot/.venv/bin/python app.py
directory=/home/oj/HWOJ/discordbot
stopsignal=QUIT
stdout_logfile=/home/oj/HWOJ/tmp/hwoj-discordbot.stdout.log
stderr_logfile=/home/oj/HWOJ/tmp/hwoj-discordbot.stderr.log